# beautybot_test.service
[Unit]
Description=BeautyBot Test Bot (Development)
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/BeautyBot_Test
Environment="PATH=/root/BeautyBot_Test/venv/bin"
ExecStart=/root/BeautyBot_Test/venv/bin/python /root/BeautyBot_Test/bot.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
